# s6/.local/share/s6/sv/pipewire-pulse/run
#!/bin/execlineb -P
exec pipewire-pulse
